# 10-keyboard-hotplug.rules
SUBSYSTEM=="input", ACTION=="add", RUN+="/bin/systemctl start --no-block autoxkbmap.service"
